# systemd/change-tty.service
# /etc/systemd/system

[Unit]
Description=Change TTY after boot

[Service]
ExecStart=/usr/bin/chvt 2

[Install]
WantedBy=default.target
