# systemd/instiki.service
[Unit]
Description=A Personal Wiki

[Service]
ExecStart=/home/user/.local/srv/instiki/instiki

[Install]
WantedBy=default.target
