# systemd/instiki.service
[Unit]
Description=A Personal Wiki

[Service]
ExecStart=/home/pavel/.webapps/instiki/instiki

[Install]
WantedBy=default.target
